# Postfix-Scripts/extract_e3k_recipients/extract.cron.in
0,15,30,45 * * * * cd /root/relay_recipients && /usr/bin/make
